# rootfs/etc/s6-overlay/s6-rc.d/mlat-client/run
#!/bin/sh
exec /etc/s6-overlay/scripts/mlat-client
