# rootfs/etc/s6-overlay/s6-rc.d/gpsd_mlat_restart/run
#!/bin/sh
exec /etc/s6-overlay/scripts/gpsd_mlat_restart
