# rootfs/etc/s6-overlay/s6-rc.d/cleanup_globe_history/run
#!/bin/sh
exec /etc/s6-overlay/scripts/cleanup_globe_history
